# supervisord.conf
[supervisord]
logfile=/root/kmd-sync-api/logs/supervisord.log
nodaemon=false
loglevel=info
user=root

[program:sync-api]
command=/usr/bin/python3 /root/kmd-sync-api/sync-api.py
stdout_logfile=/root/kmd-sync-api/logs/sync-api.out.log
stderr_logfile=/root/kmd-sync-api/logs/sync-api.err.log
autostart=true
autorestart=true
startretries=1
user=root
numprocs=1
